# nginx/crontab
0	1	*	*	*	/duckdns/duck.sh > /dev/null 2>&1
0	1	*	*	*	certbot renew
#empty str
